# mistralvibe/rootfs/etc/services.d/vibe-terminal/finish
#!/usr/bin/with-contenv bashio
if [ "${1}" -ne 0 ] && [ "${1}" -ne 256 ]; then
    bashio::log.warning "vibe-terminal exited with error code ${1}"
fi
