# mistralvibe/rootfs/etc/services.d/vibe-terminal/run
#!/usr/bin/with-contenv bashio

export VIBE_HOME=/data/vibe
if [ -f "${VIBE_HOME}/.env" ]; then
    set -a
    . "${VIBE_HOME}/.env"
    set +a
fi

mkdir -p /data/vibe/.bash
export HISTFILE=/data/vibe/.bash/history

cat > /root/.bashrc <<'BASHRC'
export VIBE_HOME=/data/vibe
[ -f "${VIBE_HOME}/.env" ] && { set -a; . "${VIBE_HOME}/.env"; set +a; }
export PS1='\[\e[1;35m\]vibe\[\e[0m\]:\[\e[1;34m\]\w\[\e[0m\]$ '
export SUPERVISOR_TOKEN="${HA_TOKEN}"
export HISTFILE=/data/vibe/.bash/history
alias vibe-restart='pkill -f mistral-vibe || true; vibe-launcher'
alias m='tmux detach 2>/dev/null || echo "Not in tmux session"'
alias v='tmux attach -t vibe 2>/dev/null || tmux new-session -A -s vibe /usr/bin/vibe-launcher'
alias ha-config='cd /config && ls'
alias ha-logs='ha core logs 2>&1 | tail -50'
alias ha-errors='ha core logs 2>&1 | grep -iE "(error|exception)" | tail -50'
BASHRC

LOG_LEVEL="$(bashio::config 'log_level')"
case "${LOG_LEVEL}" in
    trace|debug) TTYD_DEBUG="-d 7" ;;
    info|notice) TTYD_DEBUG="-d 3" ;;
    *)           TTYD_DEBUG="-d 1" ;;
esac

SESSION_PERSISTENCE="$(bashio::config 'session_persistence')"
FONT_SIZE="$(bashio::config 'terminal_font_size')"
THEME="$(bashio::config 'terminal_theme')"

THEME_JSON="$([ "${THEME:-dark}" = "light" ] && echo '{"background":"#eff1f5","foreground":"#4c4f69","cursor":"#dc8a78"}' || echo '{"background":"#0f0f17","foreground":"#e6e6f0","cursor":"#ff7a90"}')"

TTYD_OPTS="--writable --max-clients 5 --ping-interval 30 -t scrollback=20000 -t titleFixed=Mistral Vibe -t fontSize=${FONT_SIZE:-14} -t theme=${THEME_JSON} ${TTYD_DEBUG}"

if [ "${SESSION_PERSISTENCE}" = "true" ]; then
    mkdir -p /root/.config/tmux
    cat > /root/.config/tmux/tmux.conf <<'TMUXEOF'
set -g mouse on
set -g history-limit 10000
set -g default-terminal "xterm-256color"
TMUXEOF
    DESKTOP_CMD="tmux new-session -A -s vibe '/usr/bin/vibe-launcher'; bash"
else
    DESKTOP_CMD="/usr/bin/vibe-launcher"
fi

MOBILE_CMD="/usr/bin/vibe-launcher"

bashio::log.info "Starting ttyd desktop on port 7682"
ttyd --port 7682 --interface 127.0.0.1 ${TTYD_OPTS} bash -c "${DESKTOP_CMD}" &

bashio::log.info "Starting ttyd mobile on port 7683"
ttyd --port 7683 --interface 127.0.0.1 ${TTYD_OPTS} bash -c "${MOBILE_CMD}" &

bashio::log.info "Starting mobile-detection proxy on port 7681"
exec python3 /usr/share/vibe-proxy/proxy.py
